# smena.service
[Unit]
Description=Test Telegram Bot (aiogram)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/YOUR_USERNAME/smena
ExecStart=/usr/bin/python3 main.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
